# deployment/itops-backend.service
[Unit]
Description=ITOps Backend (FastAPI / Uvicorn)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/opt/itops/backend
EnvironmentFile=/opt/itops/backend/.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1

# 2 workers — safe now that the durable store is PostgreSQL (concurrent
# writers OK). Matches the 2 vCPUs on the t4g.small instance.
ExecStart=/opt/itops/venv/bin/uvicorn app.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 2 \
    --log-level info

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=itops-backend

[Install]
WantedBy=multi-user.target
